# deploy/bazpos.socket
[Unit]
Description=BazPOS gunicorn socket

[Socket]
ListenStream=/run/bazpos.sock

[Install]
WantedBy=sockets.target
